# bypass-node.service
[Unit]
Description=TON Bypass Node Service
Documentation=https://github.com/bypass-ton/node
After=network.target nss-lookup.target

[Service]
ExecStart=/usr/bin/env python3 /var/lib/bypass-node/main.py
Restart=on-failure
WorkingDirectory=/var/lib/bypass-node

[Install]
WantedBy=multi-user.target
